# os/my_init-cron/etc/service/cron/run
#!/bin/bash

if which /usr/sbin/crond &>/dev/null ; then
  exec /usr/sbin/crond -n
else
  exec /usr/sbin/cron -f
fi
